# Code/User/History/5f5868d4/jNKd.service
[Unit]
Description=Startup the system for CTF

[Service]
Type=simple
ExecStart=/root/to/your/script.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
